# systemd/sleep-lock.service
[Unit]
Description=Lock the screen on resume from suspend
